# templates/default/redis.upstart.conf.erb
#!upstart
description "Redis Server"

env USER=<%= @user %>

start on startup
stop on shutdown

respawn

exec sudo -u $USER sh -c "/usr/local/bin/redis-server \
/etc/redis/redis.conf 2>&1 >> \
<%= @log %>/redis.log"
